# ubifs-root/miwifi_rd23_firmware_99b0f_1.0.90_INT.bin/squashfs-root/usr/lib/opkg/info/nginx-ssl.list
/etc/nginx/nginx.conf
/etc/nginx/fastcgi_params
/etc/nginx/koi-win
/etc/nginx/conf.d/443.conf
/etc/nginx/conf.d/80.conf
/etc/nginx/luci-backend.conf
/etc/nginx/cert.crt
/etc/nginx/win-utf
/usr/sbin/nginx
/lib/upgrade/keep.d/nginx-ssl
/etc/nginx/koi-utf
/etc/nginx/cert.key
/etc/init.d/nginx
/etc/nginx/80-real-locations.conf
/etc/nginx/mime.types
/etc/nginx/default-locations.conf
/etc/nginx/fastcgi-proxy-tcp.conf
/etc/nginx/miwifi-webinitrd.conf
/etc/nginx/force-https-locations.conf
