# ubifs-root/miwifi_rd23_firmware_8d757_1.0.91_INT.bin/squashfs-root/etc/hotplug.d/iface/99-sp_phy_down
#!/bin/ash

if [ "$ACTION" != "ifdown" ]; then
	exit 0
fi

. /lib/miwifi/lib_sp_colls.sh

hotplug_check_down "$INTERFACE"
