# ubifs-root/miwifi_rd23_firmware_99b0f_1.0.90_INT.bin/squashfs-root/etc/uci-defaults/luci-i18n-indonesian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.id="Bahasa Indonesia"
        commit luci
EOF

