# ubifs-root/miwifi_rd23_firmware_8d757_1.0.91_INT.bin/squashfs-root/etc/hotplug.d/ntp/01-rtc
#!/bin/sh

logger -p notice -t "ntpsync" "01-rtc: run because ntp sync done"

[ -c "/dev/rtc0" ] && hwclock -w -u
