# ubifs-root/miwifi_rd23_firmware_8d757_1.0.91_INT.bin/squashfs-root/etc/hotplug.d/ntp/10-miio-ntp.sh
#!/bin/sh

logger -t duke "miio ntp hook"

miio_bind.sh -c
