# ubifs-root/miwifi_rd23_firmware_8d757_1.0.91_INT.bin/squashfs-root/etc/uci-defaults/enable_force_https.sh
#!/bin/sh
# for BSI certification, force https for all web access by default in UK region
if [ "$(bdata get CountryCode)" = "UK" ]; then
	uci set nginx.main.force_https=1
	uci commit nginx
fi
