# ubifs-root/miwifi_rd23_firmware_8d757_1.0.91_INT.bin/squashfs-root/etc/uci-defaults/luci-i18n-czech
#!/bin/sh
uci batch <<-EOF
	set luci.languages.cz=Čeština
        commit luci
EOF

