# ubifs-root/miwifi_rd23_firmware_99b0f_1.0.90_INT.bin/squashfs-root/etc/init.d/enid_cap_init
#!/bin/sh /etc/rc.common

START=22
USE_PROCD=1


start_service() {
	local netmode=$(uci -q get xiaoqiang.common.NETMODE)

	[ -z "$netmode" ] && return

	local capmode=$(uci -q get xiaoqiang.common.CAP_MODE)

	if [ "$netmode" = "whc_cap" -o "$capmode" = "ap" ]; then
		/usr/sbin/topomon_action.sh cap_init &
	fi
}
