# ubifs-root/miwifi_rd23_firmware_99b0f_1.0.90_INT.bin/squashfs-root/etc/init.d/milog
#!/bin/sh /etc/rc.common
#
# This script is used for internal application traffic control
#

START=90
STOP=89

start() {
	milog.sh -m '{"tag":"sec_sys_restart"}'
}

stop() {
	/usr/sbin/syslog_sync.sh

	return 0
}
