# ubifs-root/miwifi_rd23_firmware_153e1_1.0.49_INT.bin/squashfs-root/etc/hotplug.d/inited/02-ntp_sync
#!/bin/sh

logger -p notice -t "hotplug.d" "02-ntp_sync: run because of init done."

ntp_sync() {
    sleep 10
    ubus call wan_check reset
    return
}

ntp_sync &


