# ubifs-root/miwifi_rd23_firmware_153e1_1.0.49_INT.bin/squashfs-root/usr/lib/opkg/info/base-files.list
/usr/sbin/update_whc_ob_bssid
/usr/sbin/killpidfile
/usr/sbin/configured.sh
/etc/config/miwifi
/sbin/setssid_fact
/etc/config/luci
/etc/hotplug.d/iface/70-vpn
/etc/openwrt_version
/usr/sbin/ipfilter
/etc/banner.failsafe
/usr/sbin/wanlinkprobe
/etc/init.d/powersaving
/etc/config/cfgbackup
/etc/resolv.conf
/lib/miwifi/lib_port_map.sh
/usr/sbin/ban_admin
/usr/sbin/update_vas_conf.lua
/etc/init.d/done
/usr/sbin/do_net_diagose.lua
/var
/etc/config/traffic
/etc/config/country_specialized
/etc/config/wifiwlist
/usr/sbin/get_bdata_country
/etc/config/miled
/etc/init.d/wifiradio
/lib/preinit/00_preinit.conf
/usr/sbin/autowancheck
/etc/init.d/system
/etc/preinit
/usr/sbin/iwpriv_macfilter.lua
/etc/hotplug.d/button/51-reset
/lib/miwifi/arch/lib_arch_phy_eth.sh
/bin/reboot_safe.sh
/etc/config/devicelist
/etc/init.d/hwnat
/usr/sbin/guestwifi.sh
/etc/localtime
/usr/sbin/set_thermal_protect.sh
/usr/sbin/parentalctl.sh
/lib/upgrade/platform.sh
/usr/sbin/set_mesh_role_macfilter.sh
/sbin/smp-mt76.sh
/etc/wan_check/conned.d/20-autovpn
/etc/config/nginx
/etc/shells
/bin/flash_check.sh
/usr/sbin/set_anti_rub_network.lua
/usr/sbin/lanap_mode.sh
/lib64
/lib/preinit/70_initramfs_test
/etc/device_info
/usr/sbin/set_wps_state
/usr/lib64
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/usr/sbin/checkupgrade_cpe_cwm.lua
/usr/sbin/ip_conflict.sh
/lib/miwifi/arch/lib_arch_accel.sh
/lib/miwifi/lib_phy.sh
/etc/config/smartvpn
/lib/preinit/10_indicate_preinit
/usr/sbin/nfc.lua
/usr/sbin/luci_post_init.sh
/usr/libexec/login.sh
/usr/sbin/dhcp_apclient.sh
/usr/sbin/sysapi
/usr/sbin/check_apmode_netstatus
/usr/sbin/supervisord
/etc/config/xqled
/sbin/smp-mac80211.sh
/usr/sbin/startscene.lua
/lib/preinit/05_set_preinit_iface
/etc/config/server_mapping
/usr/sbin/gensmartdns.sh
/etc/hotplug.d/button/02-mesh
/etc/config/upnpd
/usr/sbin/logupload.lua
/etc/config/macfilter
/etc/wan_check/conned.d/check-upgrade
/usr/sbin/cr_ip_conflict.lua
/etc/config/system
/usr/sbin/check_cac_time
/usr/sbin/smp.sh
/usr/sbin/sync_guest_bssid.lua
/sbin/wps
/lib/miwifi/lib_sp_colls.sh
/etc/fstab
/etc/config/miwifi_host_mapping
/etc/config/macbind
/usr/lib/os-release
/etc/sysctl.conf
/sbin/sysupgrade
/usr/sbin/nslookupc
/usr/sbin/check_accessInternet.sh
/etc/config/power_specialized
/lib/preinit/99_10_run_init
/etc/opkg/keys/88029a1ea3d4981e
/usr/sbin/nslookupt
/lib/miwifi/arch/lib_arch_port_service.sh
/lib/preinit/98_10_mtk_failsafe_init
/usr/sbin/check_encrypt_mode.lua
/lib/miwifi/miwifi_functions.sh
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/usr/sbin/crontab_record_camera.lua
/usr/sbin/kill_plugin_process.sh
/usr/lib/lua/json.lua
/usr/sbin/setchanauto.sh
/etc/shadow
/usr/sbin/wifi_update
/etc/rc.common
/etc/hotplug.d/inited/02-ntp_sync
/sbin/setmac
/etc/uci-defaults/10_migrate-shadow
/lib/preinit/02_sysinfo
/etc/config/account
/usr/sbin/logcat
/lib/preinit/39_mount_ubi_data
/lib/preinit/10_indicate_failsafe
/sbin/phyhelper
/usr/share/xiaoqiang/xiaoqiang-defaults.txt
/usr/sbin/wifiap_mode.sh
/etc/hotplug.d/iface/91-wanip_check.sh
/sbin/nfc
/etc/iwevent.d/01-feedpush.sh
/etc/os-release
/etc/board.d/02_network
/usr/sbin/ip_change_delay.lua
/usr/sbin/sysapi.firewall
/lib/miwifi/arch/lib_arch_port_map.sh
/usr/sbin/jason.sh
/lib/upgrade/boardupgrade.sh
/usr/sbin/check_apcli_connected
/etc/config/wifiblist
/usr/sbin/iperf_script.lua
/lib/dnsmasq/lease.d/10-dhcp-mitvbox.sh
/etc/config/diag
/etc/config/firewall
/lib/preinit/80_mount_root
/lib/preinit/81_urandom_seed
/sbin/chk_sta_re
/etc/hotplug.d/iface/60-wan
/usr/sbin/ip_changed.sh
/usr/sbin/lib_wifi_util.sh
/etc/config/webfilter
/run
/etc/rc.local
/usr/sbin/set_wifi_auto_ap_mode.lua
/lib/functions/uci-defaults.sh
/usr/sbin/macfilterctl
/usr/sbin/bootinfo
/usr/sbin/firewall_cpp.sh
/sbin/getmac
/usr/sbin/record_camera_by_ip.lua
/etc/ethers
/lib/functions/network.sh
/etc/init.d/sysctl
/lib/preinit/90_mount_bind_etc
/sbin/hotplug-call
/etc/config/security
/usr/sbin/runt
/etc/sysctl.d/99-min-free-kbytes.conf
/etc/crontabs/root
/lib/preinit/06_set_rps_sock_flow
/usr/sbin/getneighbor.sh
/usr/sbin/crontab_rom.sh
/etc/hotplug.d/inited/01-ipv6
/etc/config/ddns
/lib/miwifi/miwifi_core_libs.sh
/etc/hotplug.d/iface/99-mtk-lro
/etc/hotplug.d/inited/03-restart_wifi
/lib/miwifi/lib_ap_re.sh
/etc/init.d/load_indexcache
/lib/preinit/03_preinit_do_mediatek.sh
/etc/config/country_mapping
/etc/services
/sbin/smp.sh
/etc/config/vas_user
/lib/preinit/10_check_for_mtd
/etc/init.d/miot
/lib/functions/system.sh
/sbin/port_map
/bin/flash.sh
/usr/sbin/vasinfo_fw.sh
/lib/upgrade/libupgrade.sh
/bin/board_detect
/sbin/accelctrl
/usr/sbin/dhcp_do_opt43_act.sh
/sbin/setmac_all
/etc/config/xiaoqiang
/etc/init.d/urandom_seed
/lib/preinit/39_mount_tmpfs
/etc/opkg/distfeeds.conf
/lib/preinit/31_restore_nvram
/etc/sysupgrade.conf
/usr/sbin/getCountryCode
/etc/init.d/nfc
/usr/sbin/proxy_thirdparty.sh
/usr/sbin/flashtestctl.sh
/lib/upgrade/keep.d/base-files
/lib/miwifi/arch/lib_arch_sp_colls.sh
/usr/sbin/ecos_upgrade.lua
/etc/parentalctl/parentalctl.empty
/lib/preinit/40_run_failsafe_hook
/lib/config_post_ota/add_mwan3_track
/etc/group
/usr/sbin/startscene_crontab.lua
/lib/xqled/xqled.json
/etc/config/vas
/lib/preinit/30_failsafe_wait
/etc/config/port_service
/etc/config/miio_spec
/usr/sbin/nettb2
/etc/uci-defaults/13_fix_group_user
/usr/sbin/showmacs
/bin/login.sh
/usr/sbin/config_post_ota.sh
/etc/config/deviceinfo
/etc/hotplug.d/iface/97-clear-hnat
/etc/hotplug.d/net/00-sysctl
/usr/share/xiaoqiang/xiaoqiang-ignore.txt
/etc/hosts
/etc/wan_check/conned.d/40-ntpd
/etc/uci-defaults/01_mlo-config-change.sh
/lib/functions.sh
/usr/sbin/smartvpn.sh
/etc/remove.conf
/lib/miwifi/lib_accel.sh
/etc/custom_hosts
/etc/config/network
/usr/sbin/upcrash
/etc/config/region_mapping
/lib/miwifi/lib_network.sh
/usr/sbin/check-new-firmware
/etc/smartvpn/proxy.txt
/usr/sbin/vip_device_push_act.lua
/etc/config/parentalctl
/etc/init.d/timezone
/etc/init.d/autovpn
/usr/sbin/ntpsetclock
/etc/config/port_map
/usr/sbin/dfs_led_ctrl.sh
/rom/note
/etc/inittab
/etc/openwrt_release
/sbin/hwversion
/usr/sbin/secboot_upgrade_check.sh
/usr/share/xiaoqiang/xiaoqiang_version
/usr/sbin/restore_defaults.sh
/sbin/wifi
/lib/miwifi/arch/lib_arch_network.sh
/lib/wifi/ralink_common.sh
/etc/xqDb
/usr/sbin/led_ctl
/etc/profile
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/usr/share/xiaoqiang/xiaoqiang-reserved.txt
/etc/TZ
/etc/config/misc
/etc/config/vpnlist
/etc/init.d/umount
/lib/functions/service.sh
/etc/sysctl.d/10-default.conf
/etc/passwd
/usr/sbin/check_5g_state
/lib/preinit/02_default_set_state
/usr/sbin/speed_test.lua
/etc/hotplug.d/ntp/01-rtc
/etc/mtab
/etc/config/push
/usr/sbin/dhcpc_do_opt43_act.sh
/etc/init.d/boot
/etc/config/module
/etc/config/sdkfilter
/usr/sbin/macfilter
/lib/upgrade/mmc.sh
/usr/sbin/set_timezone.sh
/usr/sbin/config_pre_ota.sh
/usr/sbin/doStatPoints
/sbin/urandom_seed
/lib/functions/preinit.sh
/lib/mediatek.sh
/usr/sbin/mipctl_public.sh
/usr/sbin/syslog_upload.lua
/usr/sbin/checkupgrade.lua
/etc/config/timezone
/bin/ipcalc.sh
/lib/firewall/ipset_redirect.loader
/etc/diag.sh
/lib/miwifi/arch/lib_arch_ap_re.sh
/usr/sbin/killpid
/etc/config/ipfilter
/etc/config/ipv6
/usr/sbin/nettb
/etc/config/backup
