# ubifs-root/miwifi_rd23_firmware_153e1_1.0.49_INT.bin/squashfs-root/usr/lib/opkg/info/miniupnpd.list
/usr/sbin/miniupnpd
/lib/upgrade/keep.d/miniupnpd
/etc/hotplug.d/iface/50-miniupnpd
/usr/share/miniupnpd/firewall.include
/etc/init.d/miniupnpd
