# ubifs-root/miwifi_rd23_firmware_153e1_1.0.49_INT.bin/squashfs-root/usr/lib/opkg/info/port_service.list
/usr/sbin/port_service
/etc/init.d/port_service
/lib/miwifi/lib_port_service.sh
