# ubifs-root/miwifi_rd23_firmware_153e1_1.0.49_INT.bin/squashfs-root/usr/lib/opkg/info/smartcontroller_c.list
/usr/sbin/smartcontroller
/usr/sbin/set_customer_wifi
/usr/sbin/set_guest_wifi
/etc/init.d/smartcontroller
