# ubifs-root/miwifi_rd23_firmware_153e1_1.0.49_INT.bin/squashfs-root/usr/lib/opkg/info/sys_tester.list
/usr/sbin/memtester
/etc/init.d/systestctl
/usr/sbin/memtestctl.helper
