# ubifs-root/miwifi_rd23_firmware_8d757_1.0.91_INT.bin/squashfs-root/usr/lib/opkg/info/firewall.list
/etc/firewall.d/firewall.basic
/usr/share/fw3/helpers.conf
/lib/firewall.sysapi.loader
/etc/init.d/firewall
/usr/sbin/dmz_bypass_ctf
/etc/security/limits.conf
/sbin/fw3
/etc/firewall.user
/lib/upgrade/keep.d/firewall
/etc/hotplug.d/iface/20-firewall
