# ubifs-root/miwifi_rd23_firmware_8d757_1.0.91_INT.bin/squashfs-root/usr/lib/opkg/info/miwifi-discovery.list
/etc/init.d/miwifi-discovery
/usr/sbin/miwifi-discovery
