# ubifs-root/miwifi_rd23_firmware_8d757_1.0.91_INT.bin/squashfs-root/usr/lib/opkg/info/sigma_daemon.list
/etc/init.d/mtk_dut
/sbin/mtk_dut
/sbin/mtk_dut_start.sh
/sbin/mtk_dut_jedi_hostapd_start.sh
