# ubifs-root/miwifi_rd23_firmware_8d757_1.0.91_INT.bin/squashfs-root/usr/lib/opkg/info/syslog-ng.list
/usr/lib/syslog-ng/libafuser.so
/usr/lib/syslog-ng/libbasicfuncs.so
/usr/lib/syslog-ng/libafprog.so
/usr/lib/syslog-ng/libcryptofuncs.so
/usr/lib/syslog-ng/libaffile.so
/usr/lib/libsyslog-ng-3.5.6.so
/etc/syslog-ng.conf
/etc/init.d/syslog-ng
/usr/lib/syslog-ng/libafsocket-notls.so
/usr/lib/syslog-ng/liblinux-kmsg-format.so
/usr/sbin/syslog-ng-ctl
/usr/lib/libsyslog-ng.so
/usr/lib/syslog-ng/libconfgen.so
/usr/sbin/syslog-ng
/usr/lib/syslog-ng/libsystem-source.so
/usr/lib/syslog-ng/libsyslogformat.so
/usr/lib/syslog-ng/libcsvparser.so
/usr/lib/syslog-ng/libsyslog-ng-crypto.so
/lib/upgrade/keep.d/syslog-ng
/usr/lib/syslog-ng/libdbparser.so
/usr/lib/syslog-ng/libafstomp.so
