# ubifs-root/miwifi_rd23_firmware_8d757_1.0.91_INT.bin/squashfs-root/usr/lib/opkg/info/xqled_v2.list
/lib/xqled/xqled.sh
/usr/sbin/xqled
/lib/xqled/xqled_common.sh
/etc/init.d/xqled
/lib/xqled/xqled_hal.sh
/lib/xqled/xqled_gpio.sh
/lib/xqled/xqled_pwm.sh
/usr/sbin/led_srv
/lib/xqled/xqled_sysfs.sh
