# ubifs-root/miwifi_rd23_firmware_153e1_1.0.49_INT.bin/squashfs-root/etc/uci-defaults/01_mlo-config-change.sh
#!/bin/sh

if [ -n "$(uci -q show misc.mld)" ]; then
	uci -q set misc.mld.ap_mlo=
	uci -q set misc.mld.hostap_mlo="5g 5gh"
	uci -q set misc.mld.bh_ap_mlo="5g 5gh"
	uci commit misc
fi
