# ubifs-root/miwifi_rd23_firmware_153e1_1.0.49_INT.bin/squashfs-root/etc/hotplug.d/iface/95-ddns
#!/bin/sh
logger -p notice -t "hotplug.d" "95-ddns: run because of $INTERFACE $ACTION"
case "$ACTION" in
ifup | ifdown)
/etc/init.d/ddns restart &
;;
esac
